# guidance-of-build/references/GOB-V1.0-reference/services/gob-monitoring.service
[Unit]
Description=GOB Monitoring Dashboard - Real-time monitoring and control for GOB multi-agent system
Documentation=file:///home/ds/GOB/monitoring/README.md
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ds
Group=ds
WorkingDirectory=/home/ds/GOB/monitoring
Environment=PATH=/home/ds/miniconda3/envs/gob/bin:/home/ds/miniconda3/bin:/usr/local/bin:/usr/bin:/bin
Environment=CONDA_PREFIX=/home/ds/miniconda3/envs/gob
Environment=PYTHONPATH=/home/ds/GOB
Environment=GOB_DIRECTORY=/home/ds/GOB
Environment=DISPLAY=:0

# Use conda environment
ExecStart=/home/ds/miniconda3/envs/gob/bin/python server.py --port 8050 --gob-dir /home/ds/GOB --no-browser

# No pre-checks - start independently

# Auto-open browser window after startup (optional)
ExecStartPost=/bin/bash -c 'sleep 10 && /home/ds/GOB/monitoring/open_monitor.sh'

# Restart settings
Restart=always
RestartSec=5
StartLimitInterval=300
StartLimitBurst=5

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gob-monitoring

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=false
ProtectSystem=strict
ReadWritePaths=/home/ds/GOB/monitoring/logs /home/ds/GOB/monitoring/data

[Install]
WantedBy=multi-user.target
